# deployments/svc/restart
#!/sbin/openrc-run

command="/bin/svc"
supervisor="supervise-daemon"
supervise_daemon_args="--stdout /dev/fd/1"
extra_started_commands="reload"

depend() {
    need net
    need localmount
}

reload() {
       ebegin "Reloading svc-daemon"
       ${command}
       eend $?
}

