# deployments/svc/restart
#!/sbin/openrc-run

command="/bin/svc"
supervisor="supervise-daemon"
supervise_daemon_args="--respawn-max 0 --verbose"
output_log="/var/log/restart.log"
error_log=${output_log}

depend() {
    use logger
    need net
    need localmount
}
